# uploader/trunk-recorder/systemd/call-recorder-uploader-drain.timer
[Unit]
Description=Retry pending Call Recorder uploads

[Timer]
OnBootSec=30
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
